# systemd/portvpn-scheduler.timer
[Unit]
Description=PortVPN Scheduler Timer
Requires=portvpn-scheduler.service

[Timer]
OnCalendar=*:*:0
AccuracySec=1s
Persistent=true

[Install]
WantedBy=timers.target
